# roles/barrybase-hdt22.py/templates/berrybase-dht22.service.j2
[Unit]
Description=Barrybase DHT22 humidity & temperature sensor
After=multi-user.target
Conflicts=getty@tty1.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 {{ install_path }}/Berrybase-DHT22.py
StandardInput=tty-force

[Install]
WantedBy=multi-user.target
